# challenges/Kurapika/deploy/Kurapika/WAF/config/supervisord.conf
[supervisord]
nodaemon=true

[program:waf]
command=go run /tmp/main.go
autostart=yes
autorestart=yes
sterr_logfile=/tmp/waf_sterr_log
stdout_logfile=/tmp/waf_stdout_log